# etc/systemd/system/evtest-grab@.service
[Service]
ExecStart="/bin/evtest --grab /dev/input/%i"
Restart=always
