# etc/systemd/system/evtest-grab@.service
[Unit]
Description=Input Grabber

[Service]
ExecStart=sh -c '/bin/evtest --grab /dev/input/%i >/dev/null 2>&1'
Restart=always
